# openboxwm.service

[Unit]
Description=Openbox Window Manager
After=dbus.socket

[Service]
ExecStart=/usr/bin/openbox

