# openbox.target
[Unit]
Description=Openbox desktop
Requires=dbus.socket
After=dbus.socket
Wants=openboxwm.service
AllowIsolate=true

